# docker-file/supervisor.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -k start

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf